# root/etc/s6-overlay/s6-rc.d/init-mod-universal-kcc-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-calibreweb-kcc-add-package/run